# rootfs/etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
pidfile=/tmp/supervisord.pid
directory=/tmp

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:arbiter]
command=/usr/bin/shinken-arbiter -c /etc/shinken/shinken.cfg
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_arbiter.log

[program:broker]
command=/usr/bin/shinken-broker -c /etc/shinken/daemons/brokerd.ini
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_broker.log

[program:carbon]
command=/opt/graphite/bin/carbon-cache.py --nodaemon start
directory=/opt/graphite
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_carbon.log

[program:graphite]
command=/usr/bin/gunicorn -b 127.0.0.1:6060 -w2 graphite_wsgi:application
directory=/opt/graphite/webapp/graphite
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_graphite.log

[program:inotify]
command=/usr/local/bin/watcher
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_inotify.log

[program:mongodb]
command=/usr/bin/mongod --dbpath /shinken/mongodb/ --logpath /tmp/mongodb.log
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_mongodb.log

[program:nginx]
command=/usr/sbin/nginx
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_nginx.log

[program:poller]
command=/usr/bin/shinken-poller -c /etc/shinken/daemons/pollerd.ini
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_poller.log

[program:reactionner]
command=/usr/bin/shinken-reactionner -c /etc/shinken/daemons/reactionnerd.ini
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_reactionner.log

[program:receiver]
command=/usr/bin/shinken-receiver -c /etc/shinken/daemons/receiverd.ini
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_receiver.log

[program:scheduler]
command=/usr/bin/shinken-scheduler -c /etc/shinken/daemons/schedulerd.ini
autostart=true
redirect_stderr=true
logfile=/tmp/supervisor_scheduler.log
